# ashnetd.service
[Unit]
Description=ashnetd

[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=ashnetd.sh

[Install]
WantedBy=multi-user.target
